# assets/etc/supervisord/ejtserver.conf
[program:ejtserver]
command = /usr/local/bin/ejtserver start-launchd
user = docker
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
